# content/usr/share/containers/systemd/zeroclaw.container
[Unit]
Description=ZeroClaw agent runtime
After=network-online.target
Wants=network-online.target

[Container]
UserNS=auto
EnvironmentFile=/etc/zeroclaw/secrets.env

# Custom image with bundled docs & source (built from /usr/share/zeroclaw/Containerfile)
Image=localhost/zeroclaw:latest
ContainerName=zeroclaw

# Publish gateway port on loopback only.
# Not exposed to LAN; access via Tailscale or SSH tunnel.
PublishPort=127.0.0.1:42617:42617

Volume=/var/opt/zeroclaw:/zeroclaw-data:z
Mount=type=bind,source=/etc/zeroclaw/secret_key,destination=/zeroclaw-data/.zeroclaw/.secret_key,ro=true

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
