# content/usr/share/containers/systemd/zeroclaw.container
[Unit]
Description=ZeroClaw agent runtime
After=network-online.target
Wants=network-online.target

[Container]
UserNS=auto
EnvironmentFile=/etc/zeroclaw/secrets.env

# Pin a release in production; :latest is distroless (no shell — use :debian to exec a shell).
Image=ghcr.io/zeroclaw-labs/zeroclaw:debian
ContainerName=zeroclaw

# Publish gateway port on loopback only.
# Not exposed to LAN; access via Tailscale or SSH tunnel.
PublishPort=127.0.0.1:42617:42617

Volume=/var/opt/zeroclaw:/zeroclaw-data:z
Mount=type=bind,source=/etc/zeroclaw/secret_key,destination=/zeroclaw-data/.zeroclaw/.secret_key,ro=true

# Optional rolling-upgrade path — re-pull a newer image on (re)start and opt into `podman auto-update`:
Pull=newer
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
